# files/app.service
[Unit]
Description=App

[Service]
ExecStart=/home/app/app

[Install]
WantedBy=multi-user.target